# target/linux/sunxi/base-files/etc/profile.d/freeioe.sh
. /lib/functions.sh
. /lib/functions/tlink.sh

export DEVICE_SN="$(device_sn)"
export IOE_CLOUD="$(ioe_cloud)"
