# target/linux/sunxi/base-files/etc/hotplug.d/tty/90-tlink
#!/bin/sh

. /lib/functions.sh

[ "$ACTION" = add ] || [ "$ACTION" = remove ] || exit 0

# Check file exists
[ -e "/dev/$DEVICENAME" ] || exit 0

board=$(board_name)

tlink_m408_xr_usb() {
	case "$DEVICENAME" in
		ttyXRUSB0)
			rm /dev/ttyS5
			ln -s /dev/$DEVICENAME /dev/ttyS5
			;;
		ttyXRUSB1)
			rm /dev/ttyS6
			ln -s /dev/$DEVICENAME /dev/ttyS6
			;;
		ttyXRUSB2)
			rm /dev/ttyS7
			ln -s /dev/$DEVICENAME /dev/ttyS7
			;;
		ttyXRUSB3)
			rm /dev/ttyS8
			ln -s /dev/$DEVICENAME /dev/ttyS8
			;;
	esac
}

tlink_m416_xr_usb() {
	case "$DEVICENAME" in
		ttyXRUSB0)
			rm /dev/ttyS5
			ln -s /dev/$DEVICENAME /dev/ttyS5
			;;
		ttyXRUSB1)
			rm /dev/ttyS6
			ln -s /dev/$DEVICENAME /dev/ttyS6
			;;
		ttyXRUSB2)
			rm /dev/ttyS7
			ln -s /dev/$DEVICENAME /dev/ttyS7
			;;
		ttyXRUSB3)
			rm /dev/ttyS8
			ln -s /dev/$DEVICENAME /dev/ttyS8
			;;
		ttyXRUSB4)
			rm /dev/ttyS13
			ln -s /dev/$DEVICENAME /dev/ttyS9
			;;
		ttyXRUSB5)
			rm /dev/ttyS14
			ln -s /dev/$DEVICENAME /dev/ttyS10
			;;
		ttyXRUSB6)
			rm /dev/ttyS15
			ln -s /dev/$DEVICENAME /dev/ttyS11
			;;
		ttyXRUSB7)
			rm /dev/ttyS16
			ln -s /dev/$DEVICENAME /dev/ttyS12
			;;
		ttyXRUSB8)
			rm /dev/ttyS13
			ln -s /dev/$DEVICENAME /dev/ttyS13
			;;
		ttyXRUSB9)
			rm /dev/ttyS14
			ln -s /dev/$DEVICENAME /dev/ttyS14
			;;
		ttyXRUSB10)
			rm /dev/ttyS15
			ln -s /dev/$DEVICENAME /dev/ttyS15
			;;
		ttyXRUSB11)
			rm /dev/ttyS16
			ln -s /dev/$DEVICENAME /dev/ttyS16
			;;
	esac
}

tlink_k1_wk_spi() {
	case "$DEVICENAME" in
		ttySWK0)
			rm /dev/ttyS7
			ln -s /dev/$DEVICENAME /dev/ttyS7
			;;
		ttySWK1)
			rm /dev/ttyS8
			ln -s /dev/$DEVICENAME /dev/ttyS8
			;;
	esac
}

tlink_k2_wk_spi() {
	case "$DEVICENAME" in
		ttySWK0)
			rm /dev/ttyS7
			ln -s /dev/$DEVICENAME /dev/ttyS7
			;;
		ttySWK1)
			rm /dev/ttyS8
			ln -s /dev/$DEVICENAME /dev/ttyS8
			;;
		ttySWK2)
			rm /dev/ttyS9
			ln -s /dev/$DEVICENAME /dev/ttyS9
			;;
		ttySWK3)
			rm /dev/ttyS10
			ln -s /dev/$DEVICENAME /dev/ttyS10
			;;
	esac
}

tlink_k2_ch432_spi() {
	case "$DEVICENAME" in
		ttyWCH0)
			rm /dev/ttyS7
			ln -s /dev/$DEVICENAME /dev/ttyS7
			;;
		ttyWCH1)
			rm /dev/ttyS8
			ln -s /dev/$DEVICENAME /dev/ttyS8
			;;
	esac
}

tlink_k4x_ch_spi() {
	case "$DEVICENAME" in
		ttyWCH0)
			rm /dev/ttyS5
			ln -s /dev/$DEVICENAME /dev/ttyS5
			;;
		ttyWCH1)
			rm /dev/ttyS6
			ln -s /dev/$DEVICENAME /dev/ttyS6
			;;
		ttyWCH2)
			rm /dev/ttyS7
			ln -s /dev/$DEVICENAME /dev/ttyS7
			;;
		ttyWCH3)
			rm /dev/ttyS8
			ln -s /dev/$DEVICENAME /dev/ttyS8
			;;
		ttyWCHB0)
			rm /dev/ttyS9
			ln -s /dev/$DEVICENAME /dev/ttyS9
			;;
		ttyWCHB1)
			rm /dev/ttyS10
			ln -s /dev/$DEVICENAME /dev/ttyS10
			;;
		ttyWCHB2)
			rm /dev/ttyS11
			ln -s /dev/$DEVICENAME /dev/ttyS11
			;;
		ttyWCHB3)
			rm /dev/ttyS12
			ln -s /dev/$DEVICENAME /dev/ttyS12
			;;
		ttyWCHC0)
			rm /dev/ttyS13
			ln -s /dev/$DEVICENAME /dev/ttyS13
			;;
		ttyWCHC1)
			rm /dev/ttyS14
			ln -s /dev/$DEVICENAME /dev/ttyS14
			;;
		ttyWCHC2)
			rm /dev/ttyS15
			ln -s /dev/$DEVICENAME /dev/ttyS15
			;;
		ttyWCHC3)
			rm /dev/ttyS16
			ln -s /dev/$DEVICENAME /dev/ttyS16
			;;
	esac
}


case $board in
kooiot,tlink-m408)
	tlink_m408_xr_usb
	;;
kooiot,tlink-m416)
	tlink_m416_xr_usb
	;;
kooiot,tlink-k1)
	tlink_k1_wk_spi
	;;
kooiot,tlink-k2|\
kooiot,tlink-k2x)
	tlink_k2_wk_spi
	tlink_k2_ch432_spi
	;;
kooiot,tlink-k4a|\
kooiot,tlink-k4b)
	tlink_k4x_ch_spi
	;;
esac

