# target/linux/sunxi/base-files/etc/hotplug.d/tty/90-xrusb
#!/bin/sh

. /lib/functions.sh

[ "$ACTION" = add ] || [ "$ACTION" = remove ] || exit 0

board=$(board_name)

tlink_m408_xr_usb() {
	case "$DEVICENAME" in
		ttyXRUSB0)
			ln -s "/dev/$DEVICENAME" "/dev/ttyS10"
			;;
		ttyXRUSB1)
			ln -s "/dev/$DEVICENAME" "/dev/ttyS11"
			;;
		ttyXRUSB2)
			ln -s "/dev/$DEVICENAME" "/dev/ttyS12"
			;;
		ttyXRUSB3)
			ln -s "/dev/$DEVICENAME" "/dev/ttyS13"
			;;
	esac
}

tlink_k1_wk_spi() {
	case "$DEVICENAME" in
		ttySWK0)
			ln -s "/dev/$DEVICENAME" "/dev/ttyS7"
			;;
		ttySWK1)
			ln -s "/dev/$DEVICENAME" "/dev/ttyS8"
			;;
	esac
}

case $board in
kooiot,tlink-m408)
	tlink_m408_xr_usb
	;;
kooiot,tlink-k1)
	tlink_k1_wk_spi
	;;
esac

