# target/linux/ramips/mt76x8/base-files/etc/hotplug.d/usb/20-cp210x
#!/bin/sh

. /lib/functions.sh

board=$(board_name)

tlink_c2_usb_serial() {
	CP210_PID="10c4/ea60/100"
	SYMLINK="ttyS2"

	if [ "${PRODUCT}" = "${CP210_PID}" ]; then
		if [ "${ACTION}" = "bind" ]; then
			DEVICE_NAME=$(ls /sys/${DEVPATH} | grep tty)
			if [ -z ${DEVICE_NAME} ]; then
				logger -t hotplug "Warning DEVICE_NAME is empty"
				exit
			fi
			logger -t hotplug "Device name of cp210 is ${DEVICE_NAME}"
			rm -f /dev/${SYMLINK}
			ln -s /dev/${DEVICE_NAME} /dev/${SYMLINK}
			logger -t hotplug "Symlink from /dev/${DEVICE_NAME} to /dev/${SYMLINK} created"
		fi

		if [ "${ACTION}" = "remove" ]; then
			rm -f /dev/${SYMLINK}
			logger -t hotplug "Symlink /dev/${SYMLINK} removed"
		fi
	fi
}

case $board in
	kooiot,tlink-c2-16m |\
	kooiot,tlink-c2-32m)
		tlink_c2_usb_serial
		;;
esac
