# target/linux/ramips/mt76x8/base-files/etc/hotplug.d/usb/20-ch341
#!/bin/sh

. /lib/functions.sh

board=$(board_name)

logger -t hotplug $(env)

tlink_c2_usb_serial() {
	CH341_PID="1a86:7523"
	SYMLINK="ttyS2"

	if [ ! -f /sys$DEVPATH/../idVendor ]; then
		return
	fi

	vid=$(cat /sys$DEVPATH/../idVendor)
	pid=$(cat /sys$DEVPATH/../idProduct)

	echo "CH341 Check""${vid}:${pid}" > /dev/kmsg

	if [ "${vid}:${pid}" = "${CH341_PID}" ]; then
		if [ "${ACTION}" = "bind" ]; then
			DEVICE_NAME=$(ls /sys/${DEVPATH} | grep tty)
			if [ -z ${DEVICE_NAME} ]; then
				logger -t hotplug "Warning DEVICE_NAME is empty"
				exit
			fi
			logger -t hotplug "Device name of ch341 is ${DEVICE_NAME}"
			rm -f /dev/${SYMLINK}
			ln -s /dev/${DEVICE_NAME} /dev/${SYMLINK}
			logger -t hotplug "Symlink from /dev/${DEVICE_NAME} to /dev/${SYMLINK} created"
		fi

		if [ "${ACTION}" = "remove" ]; then
			rm -f /dev/${SYMLINK}
			logger -t hotplug "Symlink /dev/${SYMLINK} removed"
		fi
	fi
}

case $board in
	kooiot,tlink-c2-v2)
		tlink_c2_usb_serial
		;;
esac
